# hosts/optiplex/containers/seq/seq.container
[Unit]
Description=Seq Logging Server (Containerized)
After=network-online.target
Wants=network-online.target

[Service]
TimeoutStartSec=1000
Restart=always
RestartSec=5s

[Container]
# Container image
ContainerName=seq
Image=docker.io/datalust/seq:latest
Pull=always

# Use named volume for persistent data
Volume=seq-data:/data:Z
Volume=%h/homelab/hosts/optiplex/containers/seq/seq.ini:/etc/seq/seq.ini:Z

# Network (should match your .NET app network)
Network=devops

# Environment variables for Seq
Environment=ACCEPT_EULA=Y
Environment=TZ=America/New_York
Environment=SEQ_FIRSTRUN_ADMINUSERNAME=admin
Environment=SEQ_FIRSTRUN_ADMINPASSWORD=MyStrongPassword123!
Environment=HOST_HTTP_PORT=80
Environment=HOST_INGESTION_PORT=5341

# Ports
PublishPort=8900:80
PublishPort=5341:5341
