# hosts/rpi-1/containers/postgres/patroni.container
[Unit]
Description=Patroni PostgreSQL Node (rpi-1)
After=etcd.service
Requires=etcd.service

[Container]
Image=quay.io/ongres/patroni:latest
ContainerName=patroni
Network=host
Volume=postgres-data:/var/lib/postgresql/data:Z

# Patroni identity
Environment=PATRONI_SCOPE=cluster1
Environment=PATRONI_NAME=rpi-1

# Patroni REST API
# Use 0.0.0.0 for listen and host IP for connect
Environment=PATRONI_RESTAPI_LISTEN=0.0.0.0:8008
Environment=PATRONI_RESTAPI_CONNECT_ADDRESS=$(hostname -I | awk '{print $1}'):8008

# Etcd cluster
# Use IP addresses if hostname resolution fails inside container
Environment=PATRONI_ETCD_HOSTS=$(getent hosts rpi-1 | awk '{print $1}'):2379,$(getent hosts rpi-2 | awk '{print $1}'):2379

# PostgreSQL settings
Environment=PATRONI_POSTGRESQL_DATA_DIR=/var/lib/postgresql/data
Environment=PATRONI_POSTGRESQL_LISTEN=0.0.0.0:5432
Environment=PATRONI_POSTGRESQL_CONNECT_ADDRESS=$(hostname -I | awk '{print $1}'):5432

# Replication parameters
Environment=PATRONI_POSTGRESQL_PARAMETERS_WAL_LEVEL=replica
Environment=PATRONI_POSTGRESQL_PARAMETERS_HOT_STANDBY=on
Environment=PATRONI_POSTGRESQL_PARAMETERS_MAX_WAL_SENDERS=10
Environment=PATRONI_POSTGRESQL_PARAMETERS_MAX_REPLICATION_SLOTS=10

[Service]
Restart=always
Delegate=yes
RestartSec=2
ExecStartPre=/usr/bin/chown -R 26:26 /var/lib/containers/storage/volumes/postgres-data/_data

[Install]
WantedBy=multi-user.target
